# ansible/files/exporter/labgrid-exporter.service
[Unit]
Description=Labgrid Exporter
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=0

[Service]
User={{ ansible_user }}
Group={{ ansible_user }}
Environment="PYTHONUNBUFFERED=1"
# Should contain LG_COORDINATOR configuration
EnvironmentFile=-/etc/environment
ExecStart=/home/{{ ansible_user }}/.local/bin/labgrid-exporter /etc/labgrid/exporter.yaml
Restart=always
RestartSec=30
DynamicUser=yes
# Adjust to your distribution (most often "dialout" or "tty")
SupplementaryGroups=dialout plugdev

[Install]
WantedBy=multi-user.target
